# install/roles/rustds/files/rustds.service
[Unit]
Description=rustds Service
After=network.target

[Service]
Type=simple
User=rustds
ExecStart=/home/rustds/server/startserver.sh
Restart=on-abort

[Install]
WantedBy=multi-user.target
